# unit/rtkbase_archive.timer
[Unit]
Description=Run rtkbase_archive.service everyday at 04H00

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target