# unit/rtkbase_archive.service
[Unit]
Description=RTKBase - Archiving and cleaning raw data

[Service]
Type=oneshot
User={user}
ExecStart={script_path}/archive_and_clean.sh